# rootfs/var/lib/dpkg/info/colord.list
/.
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/69-cd-sensors.rules
/lib/udev/rules.d/95-cd-devices.rules
/lib/systemd
/lib/systemd/system
/lib/systemd/system/colord.service
/usr
/usr/bin
/usr/bin/cd-create-profile
/usr/bin/cd-iccdump
/usr/bin/colormgr
/usr/bin/cd-fix-profile
/usr/lib
/usr/lib/colord
/usr/lib/colord/colord-sane
/usr/lib/colord/colord-session
/usr/lib/colord/colord
/usr/lib/i386-linux-gnu
/usr/lib/i386-linux-gnu/colord-sensors
/usr/lib/i386-linux-gnu/colord-sensors/libmunki-private.so
/usr/lib/i386-linux-gnu/colord-sensors/libcolord_sensor_colorhug.so
/usr/lib/i386-linux-gnu/colord-sensors/libcolord_sensor_argyll.so
/usr/lib/i386-linux-gnu/colord-sensors/libcolord_sensor_dummy.so
/usr/lib/i386-linux-gnu/colord-sensors/libcolord_sensor_huey.so
/usr/lib/i386-linux-gnu/colord-sensors/libhuey-private.so
/usr/lib/i386-linux-gnu/colord-sensors/libcolord_sensor_dtp94.so
/usr/lib/i386-linux-gnu/colord-sensors/libdtp94-private.so
/usr/lib/i386-linux-gnu/colord-plugins
/usr/lib/i386-linux-gnu/colord-plugins/libcd_plugin_scanner.so
/usr/lib/i386-linux-gnu/colord-plugins/libcd_plugin_sane.so
/usr/lib/i386-linux-gnu/colord-plugins/libcd_plugin_camera.so
/usr/share
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.ColorManager.service
/usr/share/dbus-1/services
/usr/share/dbus-1/services/org.freedesktop.ColorHelper.service
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/interfaces/org.freedesktop.ColorManager.xml
/usr/share/dbus-1/interfaces/org.freedesktop.ColorHelper.xml
/usr/share/dbus-1/interfaces/org.freedesktop.ColorManager.Sensor.xml
/usr/share/dbus-1/interfaces/org.freedesktop.ColorManager.Device.xml
/usr/share/dbus-1/interfaces/org.freedesktop.ColorManager.Profile.xml
/usr/share/color
/usr/share/color/icc
/usr/share/color/icc/colord
/usr/share/color/icc/colord/ProPhotoRGB.icc
/usr/share/color/icc/colord/ColorMatchRGB.icc
/usr/share/color/icc/colord/CIE-RGB.icc
/usr/share/color/icc/colord/Gamma5000K.icc
/usr/share/color/icc/colord/Crayons.icc
/usr/share/color/icc/colord/Gamma6500K.icc
/usr/share/color/icc/colord/WideGamutRGB.icc
/usr/share/color/icc/colord/AdobeRGB1998.icc
/usr/share/color/icc/colord/SwappedRedAndGreen.icc
/usr/share/color/icc/colord/PAL-RGB.icc
/usr/share/color/icc/colord/ECI-RGBv2.icc
/usr/share/color/icc/colord/AppleRGB.icc
/usr/share/color/icc/colord/NTSC-RGB.icc
/usr/share/color/icc/colord/sRGB.icc
/usr/share/color/icc/colord/BetaRGB.icc
/usr/share/color/icc/colord/EktaSpacePS5.icc
/usr/share/color/icc/colord/ECI-RGBv1.icc
/usr/share/color/icc/colord/Gamma5500K.icc
/usr/share/color/icc/colord/SMPTE-C-RGB.icc
/usr/share/color/icc/colord/DonRGB4.icc
/usr/share/color/icc/colord/Bluish.icc
/usr/share/color/icc/colord/BestRGB.icc
/usr/share/color/icc/colord/BruceRGB.icc
/usr/share/color/icc/colord/x11-colors.icc
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.color.policy
/usr/share/doc
/usr/share/doc/colord
/usr/share/doc/colord/copyright
/usr/share/glib-2.0
/usr/share/glib-2.0/schemas
/usr/share/glib-2.0/schemas/org.freedesktop.ColorHelper.gschema.xml
/usr/share/colord
/usr/share/colord/ti1
/usr/share/colord/ti1/printer-long.ti1
/usr/share/colord/ti1/display-short.ti1
/usr/share/colord/ti1/printer-normal.ti1
/usr/share/colord/ti1/printer-short.ti1
/usr/share/colord/ti1/display-long.ti1
/usr/share/colord/ti1/display-normal.ti1
/usr/share/colord/icons
/usr/share/colord/icons/spyder4-attach.svg
/usr/share/colord/icons/i1-pro-attach.svg
/usr/share/colord/icons/color-munki-photo-calibrate.svg
/usr/share/colord/icons/color-munki-photo-ambient.svg
/usr/share/colord/icons/hcfr-attach.svg
/usr/share/colord/icons/i1-monitor-attach.svg
/usr/share/colord/icons/spyder3-attach.svg
/usr/share/colord/icons/i1-display3-attach.svg
/usr/share/colord/icons/huey-attach.svg
/usr/share/colord/icons/colorhug-attach.svg
/usr/share/colord/icons/color-munki-smile-attach.svg
/usr/share/colord/icons/color-munki-photo-screen.svg
/usr/share/colord/icons/color-munki-photo-attach.svg
/usr/share/colord/icons/spyder2-attach.svg
/usr/share/colord/icons/color-munki-photo-projector.svg
/usr/share/colord/icons/dtp94-attach.svg
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/colord
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/cd-fix-profile.1.gz
/usr/share/man/man1/colormgr.1.gz
/usr/share/man/man1/colord.conf.1.gz
/usr/share/man/man1/cd-create-profile.1.gz
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.ColorManager.conf
/etc/colord.conf
/usr/share/doc/colord/changelog.Debian.gz
